# dist/ipcam-videos.service
[Unit]
Description=Cameras RTSP Recording Service

[Service]
ExecStart=/path/to/python /path/to/camera_service.py /path/to/config.json
Restart=always
TimeoutSec=100
User=username
Group=groupname
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
